# luci-theme-chuqitopd/root/etc/uci-defaults/30_luci-theme-chuqitopd
#!/bin/sh
[ -f /www/luci-static/chuqitopd/icons/ethernet.png ] && mv -f /www/luci-static/chuqitopd/icons/ethernet.png /www/luci-static/resources/icons/ethernet.png
[ -f /www/luci-static/chuqitopd/icons/ethernet_disabled.png ] && mv -f /www/luci-static/chuqitopd/icons/ethernet_disabled.png /www/luci-static/resources/icons/ethernet_disabled.png
[ -f /www/luci-static/chuqitopd/icons/port_down.png ] && mv -f /www/luci-static/chuqitopd/icons/port_down.png /www/luci-static/resources/icons/port_down.png
[ -f /www/luci-static/chuqitopd/icons/port_up.png ] && mv -f /www/luci-static/chuqitopd/icons/port_up.png /www/luci-static/resources/icons/port_up.png

uci batch <<-EOF
	set luci.themes.chuqitopd=/luci-static/chuqitopd
	set luci.main.mediaurlbase=/luci-static/chuqitopd
	commit luci
EOF
exit 0
