# luci-theme-chuqitopd/root/etc/uci-defaults/30_luci-theme-chuqitopd
#!/bin/sh
uci batch <<-EOF
	set luci.themes.chuqitopd=/luci-static/chuqitopd
	set luci.main.mediaurlbase=/luci-static/chuqitopd
	commit luci
EOF
exit 0
