# luci-app-snort/root/etc/uci-defaults/luci-app-snort
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@snort[-1]
	add ucitrack snort
	set ucitrack.@snort[-1].init=snort
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
